# usershub-service.conf
[program:usershub2]
command = APP_PATH/gunicorn_start.sh
autostart=true
autorestart=true
stdout_logfile = APP_PATH/var/log/errors_uhv2.log
redirect_stderr = true
